# tools/docker/edge/root/etc/s6-overlay/s6-rc.d/init-permissions/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# permissions
lsiown -R abc:abc \
    /opt/openems \
    /var/opt/openems
