# tools/docker/ui/root/etc/s6-overlay/s6-rc.d/init-permissions/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# permissions
lsiown -R abc:abc \
    /var/lib/nginx \
    /var/tmp/nginx \
    /etc/nginx \
    /var/log \
    /var/www/html
